# net/haproxy/pkg/haproxy.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/haproxy"
daemon_flags="-f ${HAPROXYCONF}/haproxy.cfg"

. /etc/rc.d/rc.subr

rc_cmd $1
