# mail/courier-imap/pkg/courier_imap.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/libexec/imapd.rc"

. /etc/rc.d/rc.subr

pexp=".*courierlogger.*imapd.pid"

rc_start() {
	${rcexec} "${daemon} start"
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
