# mail/courier-imap/pkg/courier_pop3.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/libexec/pop3d.rc"

. /etc/rc.d/rc.subr

pexp=".*courierlogger.*pop3d.pid"

rc_start() {
	${rcexec} "${daemon} start"
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
