# audio/hgd/pkg/hgd_playd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/hgd-playd"
daemon_user="_hgd"

. /etc/rc.d/rc.subr

rc_cmd $1
