# www/apache-httpd/pkg/httpd2.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/httpd2"

. /etc/rc.d/rc.subr

rc_pre() {
	rm -f /var/apache2/logs/{ssl_mutex,httpd.lock,accept.lock}.*
}

rc_cmd $1
