# net/miniupnp/minissdpd/pkg/minissdpd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/minissdpd"
daemon_flags="-i em0"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
