# mail/smtp-vilter/pkg/smtp_vilter.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/smtp-vilter"

. /etc/rc.d/rc.subr

pexp="smtp-vilter:"

rc_cmd $1
