# graphics/sane-backends/pkg/saned.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/libexec/saned"
daemon_flags="-a _saned"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
