# sysutils/monit/pkg/monit.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/monit"
daemon_flags="-d 300 -c ${SYSCONFDIR}/monitrc"

. /etc/rc.d/rc.subr

rc_cmd $1
