# net/icinga/core/pkg/ido2db.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/ido2db"
daemon_flags="-c ${SYSCONFDIR}/ido2db.cfg"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
