# www/sogo/pkg/sogod.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/sogod"
daemon_user="_sogo"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
