# net/ladvd/pkg/ladvd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/ladvd"
daemon_flags="-a -z"

. /etc/rc.d/rc.subr

pexp="ladvd: master"
rc_reload=NO

rc_cmd $1
