# security/sshguard/pkg/sshguard.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/sshguard"
daemon_flags="-l /var/log/authlog"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
