# net/tor/pkg/tor.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/tor"

. /etc/rc.d/rc.subr

rc_cmd $1
