# misc/memcached/pkg/memcached.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/memcached -d"
daemon_flags="-u _memcached -P /var/run/memcached.pid"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
