# net/powerdns/pkg/pdns_server.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/pdns_server --daemon=yes"
daemon_flags="--guardian=no --control-console --loglevel=9"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
