# mail/zarafa/zarafa/pkg/zarafa_search.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/zarafa-search"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
