# mail/zarafa/zarafa/pkg/zarafa_ical.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/zarafa-ical"

. /etc/rc.d/rc.subr

rc_cmd $1
