# mail/exim/pkg/exim.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/exim -bd"
daemon_flags="-q30m"

. /etc/rc.d/rc.subr

pexp="${daemon}"

rc_pre() {
	[ ! -f /etc/mailer.conf.exim ]
}

rc_cmd $1
