# comms/conserver/pkg/conserver.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/conserver -d"
daemon_flags="-n"

. /etc/rc.d/rc.subr

rc_cmd $1
