# security/samhain/pkg/yule.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/yule"

. /etc/rc.d/rc.subr

rc_cmd $1
