# misc/srcpd/pkg/srcpd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/srcpd"

. /etc/rc.d/rc.subr

rc_stop() {
	pkill -9 -f "^${pexp}"
}

rc_cmd $1
