# audio/mpd/pkg/mpd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/mpd"

. /etc/rc.d/rc.subr

rc_cmd $1
