# net/ddclient/pkg/ddclient.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/ddclient"

. /etc/rc.d/rc.subr

pexp="perl: ddclient"

rc_cmd $1
