# www/rt/pkg/rt.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/rt-server"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl -w ${daemon}"
rc_bg=YES
rc_reload=NO

rc_cmd $1
